# automatos-run/config/var/spool/cron/crontabs/root
* * * * * bash /root/script/job_instance_monitor_cronjob.sh >> /var/log/cron.log 2>&1
* * * * * echo "`date` Cron service is up and running." >> /var/log/cron.log 2>&1

